# nomad/nomad.conf
[program:nomad]
command=sudo nomad agent -dev -config /home/armand_ribouillault/nomad-discovery/nomad/config.hcl
autostart=true
autorestart=true
redirect_stderr=true
