# configs/consul.conf
[program:consul]
command=sudo consul agent -dev -ui-content-path=/ui/consul/ -config-file=/home/armrib88/nomad-discovery/configs/consul.hcl
autostart=true
autorestart=true
redirect_stderr=true
