# configs/consul.conf
[program:consul]
command=sudo consul agent -dev
autostart=true
autorestart=true
redirect_stderr=true
