# configs/nomad.conf
[program:nomad]
command=sudo nomad agent -dev
autostart=true
autorestart=true
redirect_stderr=true
