# configs/nomad.conf
[program:nomad]
command=sudo nomad agent -dev -config /home/armand_ribouillault/nomad-discovery/configs/nomad.hcl
autostart=true
autorestart=true
